# etc/supervisor/conf.d/dockerd.conf
[program:dockerd]
command=/usr/local/bin/start-dockerd

stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
